# files/etc/systemd/system/ttp229-keypad.service
# --------------------------------------------------------------------------
# Systemd service Definition for ttp229-keypad.service.
#
# The service just starts the script /usr/local/sbin/ttp229-reader.py
#
# Author: Bernhard Bablok
# License: GPL3
#
# Website: https://github.com/bablokb/pi-ttp229-keypad
#
# --------------------------------------------------------------------------

[Unit]
Description=Start TTP229 reader process
After=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/local/sbin/ttp229-reader.py

[Install]
WantedBy=multi-user.target
